# localtunnel.service
[Unit]
Description=LocalTunnel Service for ASContabilmente CRM
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html
Environment=NODE_ENV=production
ExecStart=/usr/bin/npx lt --port 80 --subdomain ascontabilemente
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=localtunnel

# Timeout per avvio/stop
TimeoutStartSec=60
TimeoutStopSec=30

# Limiti di sicurezza
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/html

[Install]
WantedBy=multi-user.target
