# node-socket.service
[Unit]
Description=Node.js Socket Service for CRM
After=network.target
After=mysql.service
Wants=mysql.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/var/www/CRM
ExecStart=/var/www/CRM/start_node.sh
ExecStop=/var/www/CRM/stop_node.sh
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/var/www/CRM/socket.pid
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Variabili d'ambiente
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
